# docker/ngnix/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /bin/bash -c "/etc/init.d/nginx start"
