# docker/tornado/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command = /bin/bash -c "/etc/init.d/redis-server start"
