# docker/mongo/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command = /bin/bash -c "/etc/init.d/mongodb start"
