# docker/mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command = /bin/bash -c "/etc/init.d/mysql start"
